# worker/scheduler/supercronic.cron
CRON_TZ=Asia/Seoul
0 19 * * 1-5 stock-report-worker
